# files/postfix/postfix-krbhelper.timer
[Unit]
Description=Refresh Postfix's kerberos cache regularly

[Timer]
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
