# files/dovecot/dovecot-krbhelper.timer
[Unit]
Description=Refresh Dovecot's kerberos cache regularly

[Timer]
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
