# files/ldap-groups-sync.timer
[Unit]
Description=Synchronize LDAP groups regularly

[Timer]
OnBootSec=1min
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
